# supervisord.ini
[program:api]
command=python3 -u /taskflow/web_server.py
autostart=true
autorestart=true
startsecs=3
startretries=10000
stdout_logfile=/var/log/supervisor/web_server.log
stderr_logfile=/var/log/supervisor/web_server_err.log

[program:checker]
command=python3 -u /taskflow/task_checker.py
directory=/taskflow
autostart=true
autorestart=true
startsecs=3
startretries=10000
stdout_logfile=/var/log/supervisor/task_checker.log
stderr_logfile=/var/log/supervisor/task_checker.log

[program:producer]
command=python3 -u /taskflow/task_producer.py
directory=/taskflow
autostart=true
autorestart=true
startsecs=3
startretries=10000
stdout_logfile=/var/log/supervisor/task_producer.log
stderr_logfile=/var/log/supervisor/task_producer.log



[program:consumer]
command=python3 -u /taskflow/task_receiver.py
directory=/taskflow
autostart=true
autorestart=true
startsecs=3
startretries=10000
stdout_logfile=/var/log/supervisor/task_receiver.log
stderr_logfile=/var/log/supervisor/task_receiver_err.log
